# systemd/postfix-watcher.service
[Unit]
Description=Postfix Watcher Service
After=network.target

[Service]
Type=simple
User=postfixwatcher
WorkingDirectory=/var/lib/postfix-watcher
ExecStart=/opt/postfix-watcher/venv/bin/postfix-watcher --config /etc/postfix-watcher.yaml
Restart=on-failure
EnvironmentFile=/etc/postfix-watcher.env

[Install]
WantedBy=multi-user.target
